# infinet-tools-backend/infinet-tools-backend.service
[Unit]
Description=InfiNet Tools Backend API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/infinet.services/infinet-tools-backend
Environment=NODE_ENV=production
Environment=PORT=3003
ExecStart=/usr/bin/node /var/www/infinet.services/infinet-tools-backend/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=infinet-tools-backend

[Install]
WantedBy=multi-user.target

